# systemd/battinfo-smart.service
[Unit]
Description=battinfo NVMe SMART cache refresh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/battinfo-smart
