# systemd/battinfo-snapshot.timer
[Unit]
Description=Daily battery health snapshot for battinfo

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15min
Unit=battinfo-snapshot.service

[Install]
WantedBy=timers.target
