# systemd/healthpanel-smart.service
[Unit]
Description=healthpanel SMART cache refresh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/healthpanel-smart
