# systemd/healthpanel-snapshot.timer
[Unit]
Description=Daily battery health snapshot for healthpanel

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15min
Unit=healthpanel-snapshot.service

[Install]
WantedBy=timers.target
